# meta-superbird/recipes-core/superbird-clock/files/superbird-clock-tick.service
[Unit]
Description=Bridgething clock tick
RequiresMountsFor=/var/lib

[Service]
Type=oneshot
ExecStart=/usr/bin/superbird-clock --save
