# meta-bridgething/recipes-core/bridgething-daemon/files/bridgething.service
[Unit]
Description=Bridgething daemon
After=network.target bridgething-init.service opt-overlay@bridgething.service
Wants=bluetooth.service
Requires=opt-overlay@bridgething.service
Conflicts=superbird-als.service
StartLimitIntervalSec=60
StartLimitBurst=3
OnFailure=bridgething-rollback.service

[Service]
Type=notify
Environment=RUST_LOG=info,bluer=warn
ExecStart=/opt/bridgething/daemon/bridgething.current
Restart=on-failure
RestartSec=2
TimeoutStartSec=30
NotifyAccess=main
WatchdogSec=0
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
