# meta-bridgething/recipes-core/bridgething-daemon/files/bridgething-opt-init.service
[Unit]
Description=Bridgething overlay first-boot seeding
After=bridgething-opt-overlay.service
Requires=bridgething-opt-overlay.service
Before=bridgething.service
ConditionPathIsMountPoint=/opt/bridgething

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/bridgething-opt-init

[Install]
WantedBy=multi-user.target
