# meta-bridgething/recipes-core/bridgething-daemon/files/bridgething-rollback.service
[Unit]
Description=Bridgething daemon rollback after restart-burst failure
After=bridgething-opt-overlay.service
Requires=bridgething-opt-overlay.service
ConditionPathIsMountPoint=/opt/bridgething

[Service]
Type=oneshot
ExecStart=/usr/libexec/bridgething-rollback
