# meta-bridgething/recipes-core/bridgething-daemon/files/bridgething-rollback.service
[Unit]
Description=Bridgething daemon rollback after restart-burst failure
After=opt-overlay@bridgething.service var-lib-bandaid.mount
Requires=var-lib-bandaid.mount

[Service]
Type=oneshot
ExecStart=/usr/libexec/bridgething-rollback
