# meta-bridgething/recipes-core/bridgething-daemon/files/bridgething.service
[Unit]
Description=Bridgething daemon (websocket + bluetooth bridge to gateways)
After=network.target bluetooth.service bridgething-init.service
Wants=bluetooth.service
Before=bridgething-chromium-kiosk.service

[Service]
Type=notify
Environment=RUST_LOG=info,bluer=warn
ExecStart=/usr/bin/bridgething
Restart=on-failure
RestartSec=2
TimeoutStartSec=30
NotifyAccess=main
WatchdogSec=0

# Daemon is a foreground process - shut down cleanly via SIGTERM and let it
# drain in-flight WS sessions + bluetooth handles.
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
