# meta-superbird/recipes-core/superbird-slot-ok/files/superbird-slot-ok.timer
[Unit]
Description=Mark A/B slot healthy a short time after boot

[Timer]
OnBootSec=60s
AccuracySec=2s

[Install]
WantedBy=timers.target
