# meta-bridgething/recipes-core/bridgething-fbpaint/bridgething-fbpaint/bridgething-fbpaint.service
[Unit]
Description=Paint 4-band test pattern at u-boot fb_addr (display handoff smoke test)
After=basic.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=no
# Wait briefly for /dev/mem to be there. With FBDEV_EMULATION off there
# is no /dev/fb0; we draw straight into u-boot's framebuffer at
# 0x1f800000 via /dev/mem mmap.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 10); do [ -c /dev/mem ] && exit 0; sleep 1; done; exit 1'
ExecStart=/usr/libexec/bridgething-fbpaint
# NOTE: backlight is intentionally left OFF on boot (pwm-backlight
# defaults to BACKLIGHT_POWER_OFF because pwm-meson clobbers the PWM
# enable bit between u-boot exit and pwm-backlight probe). To enable
# backlight from userspace:
#     echo 0 > /sys/class/backlight/backlight/bl_power
#     echo 96 > /sys/class/backlight/backlight/brightness
# To enable on boot, add `pwm-bootloader-on;` to the backlight DT node
# (the kernel patch for that property is in 0004).

[Install]
WantedBy=multi-user.target
