# meta-superbird/recipes-graphics/weston-init/files/superbird-weston.service
[Unit]
Description=Weston compositor
DefaultDependencies=no
After=seatd.service
Wants=seatd.service

[Service]
Type=simple
User=root
Group=root
Environment=XDG_RUNTIME_DIR=/run/weston
Environment=WAYLAND_DISPLAY=wayland-1
Environment=WESTON_EXTRA_FLAGS=
Environment=XCURSOR_THEME=blank
Environment=XCURSOR_PATH=/usr/share/icons
EnvironmentFile=-/etc/default/weston
ExecStartPre=/bin/mkdir -p /run/weston
ExecStartPre=/bin/chmod 700 /run/weston
ExecStartPre=/bin/sh -c 'for i in $(seq 1 200); do [ -e /dev/dri/card0 ] && exit 0; sleep 0.05; done; exit 1'
ExecStart=/bin/sh -c 'exec /usr/bin/weston --config=/etc/weston.ini --socket=wayland-1 $WESTON_EXTRA_FLAGS'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=basic.target
