# meta-superbird/recipes-core/superbird-provision/files/superbird-provision.service
[Unit]
Description=First-boot eMMC data-partition provisioning
DefaultDependencies=no
After=systemd-udevd.service
Before=local-fs-pre.target shutdown.target
Conflicts=shutdown.target
ConditionPathExists=/dev/mmcblk0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/superbird-provision
StandardError=journal+console

[Install]
WantedBy=sysinit.target
