# meta-superbird/recipes-bsp/superbird-bluetooth/files/superbird-bluetooth.service
[Unit]
Description=Superbird Bluetooth
DefaultDependencies=no
After=systemd-tmpfiles-setup.service
Before=bluetooth.service
Wants=bluetooth.service
RequiresMountsFor=/lib/firmware

[Service]
# Type=notify so bluetoothd waits for BDADDR programming before discovering hci0
Type=notify
NotifyAccess=all
ExecStart=/usr/bin/superbird-bluetooth
# force bluetoothd to re-read mgmt cache after our BDADDR write
ExecStartPost=/bin/systemctl try-restart bluetooth.service
Restart=on-failure
RestartSec=2
# cgroup cleanup so daemonized gpioset workers release REG_ON / DEV_WAKE
KillMode=control-group

[Install]
WantedBy=multi-user.target
