# meta-bridgething/recipes-core/bridgething-rotary/files/bridgething-rotary.service
[Unit]
Description=Bridgething rotary → wheel-scroll virtual mouse
After=multi-user.target
ConditionPathExists=/dev/uinput

[Service]
Type=simple
EnvironmentFile=-/etc/bridgething-rotary.conf
ExecStart=/usr/bin/bridgething-rotary
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
