# meta-bridgething/recipes-graphics/chromium-kiosk/files/bridgething-chromium-kiosk.service
[Unit]
Description=Bridgething chromium kiosk
BindsTo=bridgething-weston.service
After=bridgething-weston.service
PartOf=bridgething-weston.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/bridgething-kiosk
Environment=XDG_RUNTIME_DIR=/run/weston
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XCURSOR_THEME=blank
Environment=XCURSOR_PATH=/usr/share/icons
# block chromium's setpriority(-8) on UI/IO/GPU threads; cap drop + nice limit + nnp together
CapabilityBoundingSet=~CAP_SYS_NICE
AmbientCapabilities=
LimitNICE=20
NoNewPrivileges=yes
CacheDirectory=bridgething/chromium
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -S /run/weston/wayland-1 ] && exit 0; sleep 0.5; done; exit 1'
ExecStart=/usr/bin/bridgething-chromium-launch
Restart=on-failure
RestartSec=3
TimeoutStopSec=5

[Install]
WantedBy=bridgething-weston.service
