# meta-bridgething/recipes-graphics/chromium-kiosk/files/bridgething-chromium-kiosk.service
[Unit]
Description=Bridgething chromium kiosk (auto-launches the webapp on boot)
BindsTo=bridgething-weston.service
After=bridgething-weston.service bridgething.service
PartOf=bridgething-weston.service
# Prefer the bridgething daemon to be up before chromium navigates so the
# initial fetch hits a populated server instead of triggering chromium's
# error page. systemd treats Wants= as a soft dep when bridgething.service
# isn't installed yet, so the dev image still boots without a daemon.
Wants=bridgething.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/bridgething-kiosk
Environment=XDG_RUNTIME_DIR=/run/weston
Environment=WAYLAND_DISPLAY=wayland-1
# blank cursor theme — chromium calls libwayland-cursor on every
# pointer-enter to fetch a cursor surface, which reads XCURSOR_THEME
# from env (NOT /etc/environment.d, that's user-manager only).
# Together with weston's own cursor-theme=blank in weston.ini this
# kills every code path a cursor could be drawn through.
Environment=XCURSOR_THEME=blank
Environment=XCURSOR_PATH=/usr/share/icons
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -S /run/weston/wayland-1 ] && exit 0; sleep 0.5; done; exit 1'
ExecStart=/usr/bin/bridgething-chromium-launch
Restart=on-failure
RestartSec=3
TimeoutStopSec=5

[Install]
WantedBy=bridgething-weston.service
