# meta-bridgething/recipes-core/bridgething-clock/files/bridgething-clock.service
[Unit]
Description=Bridgething monotonic-forward clock guard
DefaultDependencies=no
Before=sysinit.target time-set.target shutdown.target
Conflicts=shutdown.target
RequiresMountsFor=/var/lib
ConditionPathIsReadWrite=/var/lib

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/bridgething-clock --restore
ExecStop=/usr/bin/bridgething-clock --save

[Install]
WantedBy=sysinit.target
