# meta-bridgething/recipes-core/bridgething-mfi-proxy/files/bridgething-mfi-proxy.service
[Unit]
Description=Bridgething MFi i2c-3 dev proxy (TCP server for cargo test)
Documentation=https://github.com/JoeyEamigh/bridgething
After=network.target
# bridgething.service owns /dev/i2c-3 during iAP2 auth and the
# backlight via the in-daemon ALS manager.
Conflicts=bridgething.service bridgething-als.service
ConditionPathExists=/dev/i2c-3
ConditionPathExists=/sys/class/backlight/backlight/brightness

[Service]
Type=simple
Environment=RUST_LOG=info

# Save current backlight brightness then blank the screen — bringup
# iteration doesn't need the kiosk visible, and the saved value is
# restored on stop/crash via ExecStopPost.
ExecStartPre=/bin/sh -c 'cat /sys/class/backlight/backlight/brightness > /run/bridgething-mfi-proxy.brightness'
ExecStartPre=/bin/sh -c 'echo 0 > /sys/class/backlight/backlight/brightness'

ExecStart=/usr/bin/bridgething-mfi-proxy 0.0.0.0:9090 /dev/i2c-3 0x10

# Always restore brightness on stop, including after a crash. /run is
# tmpfs so the saved file vanishes on reboot — that's fine, the next
# startup will resave from whatever the post-boot value is.
ExecStopPost=/bin/sh -c '\
    if [ -r /run/bridgething-mfi-proxy.brightness ]; then \
        cat /run/bridgething-mfi-proxy.brightness \
            > /sys/class/backlight/backlight/brightness; \
        rm -f /run/bridgething-mfi-proxy.brightness; \
    fi'

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
