# meta-bridgething/recipes-core/bridgething-mfi-proxy/files/bridgething-mfi-proxy.service
[Unit]
Description=Bridgething MFi i2c-3 dev proxy
Documentation=https://github.com/JoeyEamigh/bridgething
After=network.target
# bridgething.service owns /dev/i2c-3 + backlight; conflict instead of racing
Conflicts=bridgething.service superbird-als.service
ConditionPathExists=/dev/i2c-3
ConditionPathExists=/sys/class/backlight/backlight/brightness

[Service]
Type=simple
Environment=RUST_LOG=info

# save + blank backlight; ExecStopPost restores even after a crash
ExecStartPre=/bin/sh -c 'cat /sys/class/backlight/backlight/brightness > /run/bridgething-mfi-proxy.brightness'
ExecStartPre=/bin/sh -c 'echo 0 > /sys/class/backlight/backlight/brightness'

ExecStart=/usr/bin/bridgething-mfi-proxy 0.0.0.0:9090 /dev/i2c-3 0x10

ExecStopPost=/bin/sh -c '\
    if [ -r /run/bridgething-mfi-proxy.brightness ]; then \
        cat /run/bridgething-mfi-proxy.brightness \
            > /sys/class/backlight/backlight/brightness; \
        rm -f /run/bridgething-mfi-proxy.brightness; \
    fi'

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
